# supervisor/app.conf
[program:flaskapp]
command=/usr/bin/runapp.sh
autostart=true
autorestart=true
startsecs = 0
startretries=10